# files/supervisor-confd.ini
[program:confd]
stopsignal=INT
command=/scripts/start-confd.sh
redirect_stderr=true
stdout_logfile=/logs/confd.log

